# config/pureos/cron.monthly
#! /bin/bash
#
# Run at the beginning of the month via cron, out of dak's crontab.

set -e
set -o pipefail
set -u
export SCRIPTVARS=/opt/dak/config/pureos/vars
. $SCRIPTVARS

################################################################################

DATE=$(date +%Y-%m)
cd ${base}/log
touch $DATE
ln -sf $DATE current
chmod g+w $DATE
chown dak:ftpmaster $DATE
